# tp3/systemd/units/backup.timer
[Unit]
Description=lance la backup toutes les heures

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target