# tp4/scripts/mariadb/backup.service
[Unit]
Description=backup database giteadb

[Service]
Type=oneshot
ExecStart=/opt/backup_sql.sh

[Install]
WantedBy=multi-user.target